# uae-invoice.service
[Unit]
Description=UAE PINT AE E-Invoice Engine
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/uae_invoice
EnvironmentFile=/home/ubuntu/uae_invoice/.env
Environment="PATH=/home/ubuntu/uae_invoice/venv/bin"
Environment="API_KEYS=demo-key-123"
Environment="ENVIRONMENT=production"
ExecStart=/home/ubuntu/uae_invoice/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
